# sfs-to-ipset.timer
[Unit]
Description=Update IP blocklists from StopForumSpam

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target